# docker/supervisord/supervisord.conf
[supervisord]

[supervisorctl]

[program:app-uwsgi]
command = /usr/local/bin/uwsgi --ini /etc/uwsgi/uwsgi.ini --wsgi-disable-file-wrapper

[program:nginx-app]
command = /usr/sbin/nginx
